# deploy/systemd/grp-chat.service
[Unit]
# Optional — Streamlit legacy UI on :8501. React UI is the canonical frontend
# (served by nginx from /opt/grp-chat-web). Disable this unit if Streamlit not needed.
Description=GRP Chat Streamlit (legacy)
After=network.target grp-api.service

[Service]
Type=simple
User=claudeuser
WorkingDirectory=/opt/grp-chat
ExecStart=/usr/local/bin/streamlit run grp_chat.py --server.port 8501 --server.address 0.0.0.0 --server.headless true
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
