# rag-api/deploy/rag-api.service
[Unit]
Description=GRP RAG-API (multi-agent RFS analysis)
After=network.target redis-server.service elasticsearch.service ollama.service
Wants=redis-server.service elasticsearch.service ollama.service

[Service]
Type=simple
User=claudeuser
Group=claudeuser
WorkingDirectory=/opt/grp-chat/rag-api
EnvironmentFile=/etc/rag-api.env
ExecStart=/opt/grp-chat/rag-api/.venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 1 \
    --log-level info \
    --root-path /rag
Restart=on-failure
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=20

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
